# ds/consul/conf/sv/consul-server.ini
[program:consul-server]
user=foobar
command=make -C /demo/cluster start-server-fg
directory=/tmp
stdout_logfile=/var/log/consul/consul-server.log
stderr_logfile=/var/log/consul/consul-server-err.log
autostart=true
autorestart=true
startsecs=10
stopwatisecs=60
priority=8628
