# ds/consul/conf/sv/consul-node.ini
[program:consul-node]
user=foobar
command=make -C /demo/cluster start-node-fg
directory=/tmp
stdout_logfile=/var/log/consul/consul-node.log
stderr_logfile=/var/log/consul/consul-node-err.log
autostart=true
autorestart=true
startsecs=10
stopwatisecs=60
priority=8628
